# conduit-v11.service
[Unit]
Description=Conduit v11 web daemon

[Service]
ExecStart=/var/local/conduit-v11/run.sh
Restart=on-success

[Install]
WantedBy=default.target
